# monitor_dashboard.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chenrenhan.quant-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/chenrenhan/binance_bot_runtime/.venv/bin/python</string>
    <string>-u</string>
    <string>/Users/chenrenhan/binance_bot_runtime/monitor_panel.py</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8787</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/chenrenhan/binance_bot_runtime</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/chenrenhan/binance_bot_runtime/logs/monitor_panel.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chenrenhan/binance_bot_runtime/logs/monitor_panel.err.log</string>
</dict>
</plist>
